# src/sch_task_config_template.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2"
  xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2021-11-22T12:24:15</Date>
    <Author>dhi13man</Author>
    <URI>\auto_proxy</URI>
  </RegistrationInfo>
  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <IdleSettings>
      <Duration>PT10M</Duration>
      <WaitTimeout>PT1H</WaitTimeout>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
  </Settings>
  <Triggers>
    <EventTrigger>
      <StartBoundary>2021-11-22T12:24:00</StartBoundary>
      <Subscription>&lt;QueryList&gt;&lt;Query&gt;&lt;Select Path='Microsoft-Windows-WLAN-AutoConfig/Operational'&gt;*[System/EventID=8001]&lt;/Select&gt;&lt;/Query&gt;&lt;/QueryList&gt;</Subscription>
    </EventTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>{LOCATION}\main.py</Command>
      <WorkingDirectory>{LOCATION}\</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
